# cervello/vps/mycity-sentinella-dati.timer
[Unit]
Description=MyCity AD — occhi sui dati reali ogni 1 minuto (0 token; sveglia il cervello solo se scatta una soglia)

[Timer]
# VEGLIA IN TEMPO REALE (scelta di Nicola: 1 minuto). La macchina guarda i dati ogni minuto.
# È deterministico e NON usa token: il modello AI si accende SOLO quando questa sentinella
# accoda un lavoro perché una soglia è scattata (con dedup + cooldown + tetto di spesa).
# Cadenza regolabile nella finestra 1-5 min: alza OnUnitActiveSec se vuoi meno frequenza.
OnBootSec=60s
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target
