# cervello/vps/mycity-ritmo-sera.timer
[Unit]
Description=MyCity AD — report della sera ogni giorno alle 20:00

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
