# cervello/vps/mycity-giro.timer
[Unit]
Description=MyCity AD — giro di perlustrazione (DISATTIVATO: solo manuale con giro-ora.sh)

[Timer]
# DISATTIVATO: il giro automatico ogni 2h è stato spento.
# Il VPS ora serve SOLO per la chat del pannello (worker).
# Per lanciare un giro manuale: sudo bash /opt/mycity/ad-mycity/cervello/vps/giro-ora.sh
#
# Vecchia configurazione (per riattivare: decommentare e rimuovere le righe sotto):
# OnBootSec=5min
# OnUnitActiveSec=2h
# Persistent=true

# Timer inattivo: nessun trigger automatico.
# Se riattivato, partirà solo manualmente o dopo un 'systemctl start mycity-giro.service'.

[Install]
# NON abilitare: WantedBy commentato = il timer non parte al boot.
# WantedBy=timers.target
