# cervello/vps/mycity-ritmo-sera.timer
[Unit]
Description=MyCity AD — lancia il report della sera ogni giorno alle 18:00 (Piacenza)

[Timer]
OnCalendar=*-*-* 18:00:00 Europe/Rome
Persistent=true

[Install]
WantedBy=timers.target
