# cervello/vps/mycity-sentinella.service
[Unit]
Description=MyCity AD — sentinella coda lavori (auto-guarigione azioni approvate fallite)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mycity
WorkingDirectory=/opt/mycity/ad-mycity
# Fuso Piacenza: l'orario di reset quota ("resets 9:30pm (Europe/Rome)") va letto come ora italiana.
Environment=TZ=Europe/Rome
EnvironmentFile=-/opt/mycity/ad-mycity/cervello/vps/.env
ExecStart=/usr/bin/env node /opt/mycity/ad-mycity/cervello/sentinella-lavori.mjs
Nice=10

[Install]
WantedBy=multi-user.target
