# cervello/vps/mycity-giro.service
[Unit]
Description=MyCity AD — giro di perlustrazione (motore AI: Cursor 'agent' o Claude 'claude')
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mycity
WorkingDirectory=/opt/mycity/ad-mycity
EnvironmentFile=/opt/mycity/ad-mycity/cervello/vps/.env
# HOME serve al motore per trovare il login/config (~/.cursor o ~/.claude)
Environment=HOME=/home/mycity
# PATH che include la bin utente: Cursor CLI ('agent') si installa in ~/.local/bin.
Environment=PATH=/home/mycity/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Fuso di Piacenza: gli orari scritti in memoria devono essere ora-di-parete italiana (non UTC).
Environment=TZ=Europe/Rome
ExecStart=/usr/bin/env bash /opt/mycity/ad-mycity/cervello/giro.sh
