# cervello/vps/mycity-ritmo-settimana.timer
[Unit]
Description=MyCity AD — lancia la review settimanale ogni venerdì alle 15:00 (Piacenza)

[Timer]
# Ven 15:00 (spostata dalle 18:00 per non sovrapporsi al report della sera delle 18:00).
OnCalendar=Fri *-*-* 15:00:00 Europe/Rome
Persistent=true

[Install]
WantedBy=timers.target
