# cervello/vps/mycity-worker-chat.service
[Unit]
# Worker DEDICATO SOLO ALLA CHAT (Strada A, step 2). Gira ACCANTO a mycity-worker.service: il worker
# principale fa i lavori pesanti (giro/ritmo/metabolizza/azioni), questo prende SOLO le chat — così
# le tue domande nel Pannello ricevono risposta SUBITO anche mentre è in corso un giro da 45 minuti.
# Il claim atomico sulla coda impedisce che i due worker prendano lo stesso lavoro.
Description=MyCity AD — worker DEDICATO alla chat del Pannello (risposte sempre reattive)
After=network-online.target
Wants=network-online.target
# Nessun tetto ai riavvii (come il worker principale).
StartLimitIntervalSec=0

[Service]
Type=simple
User=mycity
WorkingDirectory=/opt/mycity/ad-mycity
EnvironmentFile=/opt/mycity/ad-mycity/cervello/vps/.env
Environment=HOME=/home/mycity
Environment=PATH=/home/mycity/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TZ=Europe/Rome
# LA corsia: questo processo prende solo i lavori tipo=chat.
Environment=WORKER_LANE=chat
ExecStart=/usr/bin/env bash /opt/mycity/ad-mycity/cervello/worker.sh
Restart=always
RestartSec=10
# 🧠 RECINTO DI MEMORIA (OOM 2026-07-14): il motore AI (cursor-agent, ~1.8GB virtual) ha saturato
# i 4GB del VPS e il kernel ha ucciso il processo a metà risposta → chat bloccata in_corso per
# sempre (episodi 00:06 e 01:13). Col recinto: sopra MemoryHigh il kernel rallenta/reclama solo
# questo servizio; a MemoryMax uccide DENTRO il recinto (il motore, il processo più grosso) senza
# toccare il resto del sistema — e Restart=always rialza il servizio da solo.
MemoryHigh=1G
MemoryMax=1536M
# Se il kernel uccide un figlio (il motore) dentro il recinto, systemd riavvia il servizio intero:
# meglio un worker fresco che uno zombie con il motore morto sotto.
OOMPolicy=kill
# Watchdog: come il worker principale, systemd riavvia se il loop si blocca. La chat non fa lavori
# lunghi come il giro, quindi un margine più stretto va bene, ma teniamo 3300s per semplicità/uniformità.
WatchdogSec=3300
NotifyAccess=all

[Install]
WantedBy=multi-user.target
