# cervello/vps/mycity-watch-main.timer
[Unit]
Description=MyCity AD — poll origin/main ogni 5 minuti

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
